# supervisord-rsyslogd.conf
[program:sshd]
command=/usr/sbin/rsyslogd
